# debian/cherokee.conffiles
/etc/cherokee/cherokee.conf
/etc/cherokee/icons.conf
/etc/init.d/cherokee
/etc/cherokee/mime.conf
/etc/logrotate.d/cherokee
